# src/packages/fff/fff-hoods/files/etc/uci-defaults/93-fff-hoodfile
. /lib/functions/fff/keyxchange

ln -s "$hoodfilewww" /www/hood/keyxchangev2data
ln -s "$sectorfilewww" /www/hood/sectorfile

exit 0
