# src/packages/fff/fff-network/files/etc/uci-defaults/22a-config-ports
#!/bin/sh
# Copyright 2019 Adrian Schmutzler
# License GPLv3

. /lib/functions.sh
. /lib/functions/fff/networksetup

rename_vlan() {
	local name="$1"
	local vlan="$(uci get network.$name.vlan)"
	uci rename network.$name="vlan$vlan"
}

BOARD="$(uci get board.model.name)"

# force config recreation (since old config might have been preserved)
rm -f /etc/config/network
/bin/config_generate

# set up ethmesh
uci set network.ethmesh=interface
uci set network.ethmesh.proto=batadv
uci set network.ethmesh.mesh=bat0

# set up mesh
uci set network.mesh=interface
uci set network.mesh.type=bridge
uci set network.mesh.auto=1

# set up mesh member interfaces
MESHIF="$(uci -q get network.lan.ifname)"
if [ -n "$MESHIF" ]; then
	uci set network.mesh.ifname="bat0 $MESHIF"
else
	uci set network.mesh.ifname="bat0"
fi

# set up vlans
config_load network
config_foreach rename_vlan switch_vlan

SWITCHHW="$(uci -q get network.@switch[0].name)"
if [ -n "$SWITCHHW" ]; then
	uci set "network.vlan3=switch_vlan"
	uci set "network.vlan3.device=$SWITCHHW"
	uci set "network.vlan3.vlan=3"
fi

case "$BOARD" in
	archer-c7-v2)
		setupSwitch "4 5 0t" "2 3 0t" "1 6t"
		;;
	archer-c7-v4|\
	archer-c7-v5|\
	tl-wdr3600|\
	tl-wdr4300|\
	tl-wdr4310|\
	tl-wdr3600-v1|\
	tl-wdr4300-v1|\
	tl-wdr4310-v1|\
	tl-wdr4900-v1)
		setupSwitch "4 5 0t" "2 3 0t" "1 0t"
		;;
	archer-c25-v1|\
	archer-c60-v1|\
	archer-c60-v2|\
	tl-wdr3500|\
	tl-wdr3500-v1|\
	tl-wr741nd-v2|\
	tl-wr841-v7|\
	tl-wr841-v9|\
	tl-wr841-v10|\
	tl-wr841-v11|\
	tl-wr841-v12)
		setupSwitch "1 2 0t" "3 4 0t"
		;;
	c2600)
		setupSwitch "1 2 6t" "3 4 6t" "5 0t"
		;;
	cpe210|\
	cpe510)
		# Default: LAN0: WAN, LAN1: CLIENT
		setupSwitch "0t 4" "0t" "0t 5"
		uci set "fff.ui.portsetup=cpev1"
		;;
	fritzbox-4040)
		# Set whole switch to CLIENT, since VLANs are not supported
		# (eth0 is connected untagged to switch)
		uci set "fff.ui.portsetup=twoeth"
		;;
	gl-ar150)
		# Default for LAN-Port: CLIENT
		setupSwitch "0t 1" "0t"
		uci set "fff.ui.portsetup=2ndport"
		;;
	tl-wr740n-v4|\
	tl-wr740nd-v4|\
	tl-wr741nd-v4|\
	tl-wr841-v8|\
	tl-wr842n-v2)
		setupSwitch "1 4 0t" "2 3 0t"
		;;
	tl-wr1043nd-v1)
		setupSwitch "3 4 5t" "1 2 5t" "0 5t"
		;;
	tl-wr1043nd-v2|\
	tl-wr1043nd-v3)
		setupSwitch "1 2 6t" "3 4 6t" "5 6t"
		;;
	tl-wr1043nd-v4|\
	tl-wr1043n-v5)
		setupSwitch "1 2 0t" "3 4 0t" "5 0t"
		;;
	ubnt-nano-m|\
	ubnt-nano-m-xw|\
	ubnt-nanostation-m|\
	ubnt-nanostation-m-xw)
		# Default: LAN0: WAN, LAN1: CLIENT
		uci set "fff.ui.portsetup=twoeth"
		;;
	cpe210-v2|\
	cpe210-v3|\
	tl-mr3020-v1|\
	tl-wa850re-v1|\
	tl-wa860re-v1|\
	tl-wa901nd-v2|\
	ubnt-bullet-m|\
	ubnt-bullet-m-xw|\
	ubnt-loco-m|\
	ubnt-loco-m-xw|\
	ubnt-pico-m|\
	ubnt-power-m-xw|\
	ubnt-unifi|\
	ubnt-unifiac-lite|\
	ubnt-unifiac-mesh)
		uci set "fff.ui.portsetup=oneport"
		;;
esac

# delete unused OpenWrt nodes
uci del network.lan # only delete AFTER switch setup
# uci del network.wan_dev # not deleting this might preserve MAC address
uci -q del network.wan6

# setup sysctl conf for WAN
WANDEV=$(uci -q get network.wan.ifname)
[ -n "$WANDEV" ] && enableAutoConf "$WANDEV"

uci commit network
uci commit fff
