# src/packages/fff/fff-network/files/etc/uci-defaults/22b-config-routermac
#!/bin/sh
# Copyright 2019 Adrian Schmutzler
# License GPLv3

BOARD="$(uci get board.model.name)"

. /lib/functions/fff/networksetup

. /lib/functions.sh
. /lib/functions/system.sh

# Set MAC address for br-mesh
case "$BOARD" in
	archer-c7-v2|\
	c2600)
		ROUTERMAC=$(cat /sys/class/net/eth1/address)
		;;
	archer-c7-v4|\
	archer-c7-v5|\
	fritzbox-4040|\
	ubnt-unifiac-lite|\
	ubnt-unifiac-mesh)
		ROUTERMAC=$(cat /sys/class/net/eth0/address)
		;;
	archer-c25-v1|\
	archer-c60-v1|\
	archer-c60-v2|\
	tl-wdr3600|\
	tl-wdr4300|\
	tl-wdr4310|\
	tl-wdr3600-v1|\
	tl-wdr4300-v1|\
	tl-wdr4310-v1)
		ROUTERMAC=$(cat /sys/class/ieee80211/phy1/macaddress)
		;;
	*)
		ROUTERMAC=$(cat /sys/class/ieee80211/phy0/macaddress)
		;;
esac

echo "Setting br-mesh MAC address to $ROUTERMAC"
uci set "network.mesh.macaddr=$ROUTERMAC"

# Set MAC address for ethmesh
case "$BOARD" in
	archer-c7-v4|\
	archer-c7-v5|\
	c2600|\
	tl-wdr4900-v1|\
	tl-wr841-v9|\
	tl-wr841-v10|\
	tl-wr841-v11|\
	tl-wr841-v12|\
	tl-wr1043nd-v4|\
	tl-wr1043n-v5)
		ETHMESHMAC=$(macaddr_add "$ROUTERMAC" 1)
		;;
	cpe210|\
	cpe210-v2|\
	cpe210-v3|\
	cpe510|\
	gl-ar150|\
	tl-mr3020-v1|\
	tl-wa850re-v1|\
	tl-wa860re-v1|\
	tl-wa901nd-v2|\
	tl-wr1043nd-v1)
		ETHMESHMAC=$(macaddr_setbit_la "$ROUTERMAC")
		;;
	tl-wr841-v8|\
	tl-wr842n-v2)
		ETHMESHMAC=$(macaddr_add "$ROUTERMAC" -1)
		;;
esac

if [ -n "$ETHMESHMAC" ]; then
	echo "Setting ethmesh MAC address to $ETHMESHMAC"

	ETHMESHDEV="$(uci -q get network.ethmesh.ifname)"
	[ -n "$ETHMESHDEV" ] || ETHMESHDEV=eth0
	uci set network.ethmesh_dev=device
	uci set network.ethmesh_dev.name="$ETHMESHDEV"
	uci set network.ethmesh_dev.macaddr="$ETHMESHMAC"
fi

uci commit network
