# src/packages/fff/fff-network/files/etc/uci-defaults/90-firstbootfff
#!/bin/sh

# process first-boot user commands
[ -s "/etc/firstbootfff" ] && {
	sh /etc/firstbootfff
}
