# src/packages/fff/fff-support/files/etc/uci-defaults/90-firstbootfff
#!/bin/sh
# Copyright 2019 Adrian Schmutzler
# License GPLv3

LOCALUPGRADESCRIPT="/etc/firstbootfff"

#add LOCALUPGRADESCRIPT to sysupgrade
grep -q ^"${LOCALUPGRADESCRIPT}"$ /etc/sysupgrade.conf || echo "$LOCALUPGRADESCRIPT" >> /etc/sysupgrade.conf

# process user commands
[ -s "$LOCALCONFIGSCRIPT" ] && {
	sh /etc/firstbootfff
}
