# src/packages/fff/fff-support/files/etc/uci-defaults/96-fff-rclocaluserconfig
echo "/etc/rc.local.fff_userconfig" >> /etc/sysupgrade.conf

exit 0
