# src/packages/fff/fff-gateway/files/etc/uci-defaults/93-fff-hoodfile
# Migrate old hood file
[ -s /www/hood/keyxchangev2data ] && mv /www/hood/keyxchangev2data /etc/hoodfile

# Link to new location
ln -s /etc/hoodfile /www/hood/keyxchangev2data

exit 0
