# src/packages/fff/fff-gateway/files/etc/uci-defaults/96-fff-gwupgrade
#!/bin/sh

echo "/etc/config/gateway" >> /etc/sysupgrade.conf
