# src/packages/fff/fff-gateway/files/etc/uci-defaults/96-fff-gwupgrade
echo "/etc/config/gateway" >> /etc/sysupgrade.conf

exit 0
