# src/packages/fff/fff-web-ui/files/etc/uci-defaults/06-fff-reset-update-notification
if [ "$(uci -q get fff.notifyupdate.resetonupgrade)" = "1" ] ; then
	uci set "fff.notifyupdate.value=1"
	uci commit fff
	/bin/rm -f /tmp/fwcheck
fi

exit 0
