# src/packages/fff/fff-jubtl/files/etc/uci-defaults/91-jubtl-public-key
if [ ! -s /etc/dropbear/authorized_keys ] ; then
	cat /etc/id_jubtl.pub >> /etc/dropbear/authorized_keys
fi
chmod 700 /etc/dropbear
chmod 600 /etc/dropbear/authorized_keys

exit 0
