# src/packages/fff/fff-jubtl/files/etc/uci-defaults/94-fff-uhttpd-cert
#!/bin/sh

uci set uhttpd.ssl.cert='/etc/ffsslchain.pem'
uci set uhttpd.ssl.key='/etc/ffssl.key'
uci commit uhttpd

exit 0
