# src/packages/fff/fff-started/files/etc/init.d/zdone
#!/bin/sh /etc/rc.common

# Copyright 2017 Adrian Schmutzler
# License GPLv3

START=99

start()
{
	touch /tmp/started
}
