# src/packages/fff/fff-wireless/files/etc/uci-defaults/24b-config-wireless
#!/bin/sh
# Copyright 2018 Adrian Schmutzler
# License GPLv3

BOARD="$(uci get board.model.name)"

# antenna gain
case "$BOARD" in
	*tl-wr1043nd-v1)
		uci set wireless.radio0.antenna_gain=3
		uci commit wireless
		;;
	*loco-m|\
	*loco-m-xw)
		uci set wireless.radio0.antenna_gain=9
		uci commit wireless
		;;
	*nanostation-m|\
	*nanostation-m-xw)
		uci set wireless.radio0.antenna_gain=11
		uci commit wireless
		;;
	pico-m)
		uci set wireless.radio0.antenna_gain=12
		uci commit wireless
		;;
	power-m-xw)
		uci set wireless.radio0.antenna_gain=18
		uci commit wireless
		;;
	*unifi)
		# No support for Unifi AP Pro!
		uci set wireless.radio0.antenna_gain=3
		uci commit wireless
		;;
	*unifiac-mesh)
		uci set wireless.radio0.antenna_gain=4
		uci set wireless.radio1.antenna_gain=3
		uci commit wireless
		;;
	*unifiac-mesh-pro)
		uci set wireless.radio0.antenna_gain=8
		uci set wireless.radio1.antenna_gain=8
		uci commit wireless
		;;
	*unifiac-pro)
		uci set wireless.radio0.antenna_gain=6
		uci set wireless.radio1.antenna_gain=3
		uci commit wireless
		;;
esac

# rssileds: always disable if package is installed
[ -x "/etc/init.d/rssileds" ] && /etc/init.d/rssileds disable
