# src/packages/fff/fff-wireless/files/etc/uci-defaults/24b-config-wireless
#!/bin/sh
# Copyright 2018 Adrian Schmutzler
# License GPLv3

BOARD="$(uci get board.model.name)"

# antenna gain
case "$BOARD" in
	tl-wr1043nd-v1)
		uci set wireless.radio0.antenna_gain=3
		uci commit wireless
		;;
	ubnt-loco-m|\
	ubnt-loco-m-xw|\
	ubnt_loco-m|\
	ubnt_loco-m-xw)
		uci set wireless.radio0.antenna_gain=9
		uci commit wireless
		;;
	ubnt-nano-m|\
	ubnt-nano-m-xw|\
	ubnt_nanostation-m|\
	ubnt_nanostation-m-xw)
		uci set wireless.radio0.antenna_gain=11
		uci commit wireless
		;;
	ubnt-pico-m|\
	ubnt_pico-m)
		uci set wireless.radio0.antenna_gain=12
		uci commit wireless
		;;
	ubnt-power-m-xw|\
	ubnt_power-m-xw)
		uci set wireless.radio0.antenna_gain=18
		uci commit wireless
		;;
	ubnt-unifi|\
	ubnt_unifi)
		# No support for Unifi AP Pro!
		uci set wireless.radio0.antenna_gain=3
		uci commit wireless
		;;
	ubnt-unifiac-mesh|\
	ubnt_unifiac-mesh)
		uci set wireless.radio0.antenna_gain=4
		uci set wireless.radio1.antenna_gain=3
		uci commit wireless
		;;
esac

# rssid
case "$BOARD" in
	cpe210|\
	cpe210-v2|\
	cpe210-v3|\
	cpe510|\
	tl-wa850re-v1)
		uci set system.rssid_wlan0.dev=radio0
		uci commit system
		;;
esac
